# supervisord/celery_grobiz.ini
[program:celery_growbiz]
user=root
environment = 
    export C_FORCE_ROOT="true"

directory=/var/www/html/prismerp/growbiz_prismerp_global/divineba
command= /usr/bin/python /var/www/html/prismerp/growbiz_prismerp_global/divineba/manage.py celeryd --verbosity=2 --loglevel=info --pool=solo --beat
 
autostart=true
autorestart=true
stdout_logfile=/var/log/celery/growbiz_prismerp_global.log
stderr_logfile=/var/log/celery/growbiz_prismerp_global.err.log